# share/cutback.timer
[Unit]
Description=Daily system backup

[Timer]
OnCalendar=daily

