# share/cutback.service
[Unit]
Description=Perform system backup

[Service]
Type=oneshot
ExecStart=/usr/bin/cutback -c /etc/xdg/cutback/config.yml

