# share/cutback.service
[Unit]
Description=Perform system backup

[Service]
Type=oneshot
ExecStart=/usr/bin/cutback -c /mnt/backup/cutback.yml

